# systemd/novahiz-mcp-server.service
[Unit]
Description=Novahiz OS MCP Server (STDIO mode)
After=network.target novahiz-runtime.service novahiz-bridge.service

[Service]
Type=simple
User=novahiz
WorkingDirectory=/home/novahiz/.opencode
EnvironmentFile=/home/novahiz/.novahiz/.env
ExecStart=/usr/bin/python3 /home/novahiz/.opencode/mcp/novahiz-mcp.py --mcp
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
